# release/aimakeimage/deploy/aimakeimage.service
[Unit]
Description=AIMakeImage Spring Boot API
After=network.target

[Service]
Type=simple
User=aimakeimage
Group=aimakeimage
WorkingDirectory=/opt/aimakeimage/backend
EnvironmentFile=/etc/aimakeimage/aimakeimage.env
ExecStart=/usr/bin/java -jar /opt/aimakeimage/backend/make-image-api.jar
Restart=always
RestartSec=5
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
